# backup/crontab
#daily at 3 am
0 3 * * * /app/backup.sh
